# systemd/.config/systemd/user/megasync.service
[Unit]
Description=MEGAsync

[Service]
Type=simple
ExecStart= /usr/bin/megasync
Restart=always

[Install]
WantedBy=graphical.target
